# .config/etc/systemd/system/getty@.service.d/keyboard-layout.conf
[Service]
ExecStartPre=/usr/bin/loadkeys /home/jean/.config/keyboard/vc.conf
